# systemd/hardn-monitor.service
[Unit]
Description=HARDN Centralized Monitoring and Logging Service
# Bring monitor up after network + core services so it can attach to logs/metrics
After=network-online.target hardn.service hardn-api.service legion-daemon.service
Wants=network-online.target hardn.service hardn-api.service legion-daemon.service
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/usr/share/hardn

ExecStart=/usr/bin/hardn-monitor

# Only restart on failure so manual stops don't immediately respawn
Restart=on-failure
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=hardn-monitor

CPUAccounting=true
IOAccounting=true
MemoryAccounting=true

# Security settings
NoNewPrivileges=yes
ProtectHome=yes
ProtectSystem=strict
# Allow writes where the monitor stores state / logs
ReadWritePaths=/var/lib/hardn /var/log/hardn
PrivateTmp=yes

# Resource limits
MemoryMax=256M
CPUQuota=50%

[Install]
WantedBy=multi-user.target