# systemd/hardn.service
[Unit]
Description=HARDN Security Monitoring and Response Service
After=network.target hardn-monitor.service
Wants=network.target hardn-monitor.service
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=root
Group=root
ExecStart=/bin/bash -c 'if [ ! -d /var/lib/hardn/legion ] || [ -z "$(ls -A /var/lib/hardn/legion/baseline_*.json 2>/dev/null)" ]; then echo "Creating LEGION baseline..."; /usr/bin/hardn legion --create-baseline --json; fi; exec /usr/bin/hardn legion --daemon --response-enabled --verbose'
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hardn
CPUAccounting=true
IOAccounting=true
MemoryAccounting=true

# Security settings
NoNewPrivileges=yes
ProtectHome=yes
ProtectSystem=strict
ReadWritePaths=/var/lib/hardn /var/log/hardn
PrivateTmp=yes

# Resource limits
MemoryMax=512M
CPUQuota=75%

[Install]
WantedBy=multi-user.target