# systemd/hardn-api.service
[Unit]
Description=HARDN API Server
# Bring API up after network + core daemon, but not after monitor
After=network-online.target hardn.service
Wants=network-online.target hardn.service
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=hardn
Group=hardn
WorkingDirectory=/usr/share/hardn

# -u / PYTHONUNBUFFERED for real-time logs in journald
Environment=PYTHONUNBUFFERED=1
# Use env to tolerate python3 path differences across distros after seeing issue between deb 13 and ubuntu 24.04
ExecStart=/usr/bin/env python3 -u /usr/share/hardn/src/hardn-api.py

Restart=on-failure
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=hardn-api

CPUAccounting=true
IOAccounting=true
MemoryAccounting=true

# Security settings
NoNewPrivileges=yes
ProtectHome=yes
ProtectSystem=strict

# Allow writes where the API is expected to log / persist
ReadWritePaths=/var/log/hardn /var/lib/hardn

PrivateTmp=yes

[Install]
WantedBy=multi-user.target