# systemd/hardn-grafana.service
[Unit]
Description=HARDN Grafana Visualization Service
After=network.target hardn.service hardn-monitor.service
Wants=network.target hardn.service hardn-monitor.service

[Service]
Type=simple
User=grafana
Group=grafana
Environment=GF_PATHS_DATA=/var/lib/grafana
Environment=GF_PATHS_LOGS=/var/log/grafana
Environment=GF_PATHS_PLUGINS=/var/lib/grafana/plugins
Environment=GF_SERVER_HTTP_PORT=3000
ExecStart=/usr/sbin/grafana-server --config=/etc/grafana/grafana.ini --homepath=/usr/share/grafana
Restart=always
RestartSec=10
LimitNOFILE=65536
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hardn-grafana

# Security Hardening
NoNewPrivileges=yes
ProtectHome=yes
ProtectSystem=strict
ReadWritePaths=/var/lib/grafana /var/log/grafana
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
