# systemd/hardn-api.service
[Unit]
Description=HARDN API Server
# Bring API up after network + core daemon, but not after monitor
After=network-online.target hardn.service
Wants=network-online.target hardn.service
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=hardn
Group=hardn
WorkingDirectory=/usr/share/hardn

# -u / PYTHONUNBUFFERED for real-time logs in journald
Environment=PYTHONUNBUFFERED=1
# Default bind is loopback only. To expose the API to other hosts set both
# HARDN_API_HOST=0.0.0.0 AND HARDN_API_ALLOWED_CIDRS=10.0.0.0/8,...; the
# Python middleware rejects any request whose source IP is not inside the
# allowlist, so a misconfigured bind without a matching CIDR list still
# fails closed.
Environment=HARDN_API_HOST=127.0.0.1
Environment=HARDN_API_PORT=8000
Environment=HARDN_API_ALLOWED_CIDRS=127.0.0.0/8,::1/128
# Authorized SSH public keys for remote API auth, one key per line
Environment=HARDN_AUTHORIZED_KEYS=/etc/hardn/authorized_keys
# Refuse to start the API if no authorized keys are registered. Even with
# the bind defaulting to 127.0.0.1, a local process can still hit the API
# and we want bearer-auth to mean something. Operators must register at
# least one key:
#   echo "$(cat ~/.ssh/id_ed25519.pub)" | sudo tee -a /etc/hardn/authorized_keys
ExecStartPre=/bin/sh -c 'test -s /etc/hardn/authorized_keys || { echo "HARDN API refusing to start: /etc/hardn/authorized_keys is missing or empty. See README \"Remote Access\"." >&2; exit 1; }'
# Use env to tolerate python3 path differences across distros after seeing issue between deb 13 and ubuntu 24.04
ExecStart=/usr/bin/env python3 -u /usr/share/hardn/src/hardn-api.py

Restart=on-failure
RestartSec=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=hardn-api

CPUAccounting=true
IOAccounting=true
MemoryAccounting=true

# Security settings
NoNewPrivileges=yes
ProtectHome=yes
ProtectSystem=strict

# Allow writes where the API is expected to log / persist
ReadWritePaths=/var/log/hardn /var/lib/hardn

PrivateTmp=yes

[Install]
WantedBy=multi-user.target