# systemd/legion-daemon.service
[Unit]
Description=HARDN LEGION Security Monitoring Daemon
# Bring LEGION up after the network is actually online
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=root
Group=root

# Keep everything consistent with the rest of the stack tomorevent the previous order loops 
WorkingDirectory=/usr/share/hardn

# Create necessary directories before starting
ExecStartPre=/bin/bash -c 'mkdir -p /var/lib/hardn/legion /var/log/hardn && chmod 755 /var/lib/hardn /var/lib/hardn/legion'

# Use the correct binary path and create baseline if needed
# (If hardn.service already created a baseline, this condition is a no-op.)
ExecStart=/bin/bash -c 'if [ ! -d /var/lib/hardn/legion ] || [ -z "$(ls -A /var/lib/hardn/legion/baseline_*.json 2>/dev/null)" ]; then echo "Creating LEGION baseline..."; /usr/bin/hardn legion --create-baseline --json; fi; exec /usr/bin/hardn legion --daemon --verbose'

# If we want stop to really stop it at some point, switch to on-failure
Restart=always
RestartSec=10

StandardOutput=journal
StandardError=journal
SyslogIdentifier=legion-daemon

CPUAccounting=true
IOAccounting=true
MemoryAccounting=true
# Set environment to use correct binary
Environment="PATH=/usr/bin:/bin:/usr/sbin:/sbin"

# Security settings - adjusted to avoid namespace issues
NoNewPrivileges=yes
ProtectHome=yes
ProtectSystem=strict
ReadWritePaths=/var/lib/hardn /var/log/hardn
PrivateTmp=yes

# Resource limits
MemoryMax=256M
CPUQuota=30%

[Install]
WantedBy=multi-user.target