# debian/postinst
#!/bin/bash
set -e

LEGACY_MODULE_DIR="/usr/lib/hardn/src/setup/modules"
LEGACY_TOOL_DIR="/usr/lib/hardn/src/setup/tools"

# hardn group sanity check 
ensure_hardn_account() {
    if ! getent group hardn >/dev/null; then
        if command -v addgroup >/dev/null 2>&1; then
            addgroup --system hardn || true
        elif command -v groupadd >/dev/null 2>&1; then
            groupadd -r hardn || true
        fi
    fi

    if ! getent passwd hardn >/dev/null; then
        if command -v adduser >/dev/null 2>&1; then
            adduser --system --ingroup hardn --home /var/lib/hardn --shell /usr/sbin/nologin hardn || true
        elif command -v useradd >/dev/null 2>&1; then
            useradd -r -g hardn -d /var/lib/hardn -s /usr/sbin/nologin hardn || true
        fi
    fi
}

ensure_hardn_paths() {
    mkdir -p /var/log/hardn /var/lib/hardn /var/lib/hardn/backups /var/lib/hardn/baselines 2>/dev/null || true
    chown root:hardn /var/log/hardn /var/lib/hardn /var/lib/hardn/backups /var/lib/hardn/baselines 2>/dev/null || true
    chmod 2770 /var/log/hardn /var/lib/hardn /var/lib/hardn/backups /var/lib/hardn/baselines 2>/dev/null || true

    if [ -d /usr/share/hardn ]; then
        chown -R root:hardn /usr/share/hardn 2>/dev/null || true
        find /usr/share/hardn -type d -exec chmod 2755 {} + 2>/dev/null || true
    fi

    mkdir -p /usr/share/hardn/src 2>/dev/null || true
    chown root:hardn /usr/share/hardn/src 2>/dev/null || true
    chmod 2755 /usr/share/hardn/src 2>/dev/null || true
}
        if [ -d "systemd" ]; then
            for service_file in systemd/*.service; do
                if [ -f "$service_file" ]; then
                    install -m 644 "$service_file" /lib/systemd/system/ || true
                fi
            done
        fi
        
        # Reload systemd if any services were installed + group add
        if [ -d "systemd" ] && command -v systemctl >/dev/null; then
            systemctl daemon-reload || true
        fi

RUNTIME_ROOT="/usr/share/hardn"
MODULE_DIR="${RUNTIME_ROOT}/modules"
TOOL_DIR="${RUNTIME_ROOT}/tools"
TEMPLATE_DIR="${RUNTIME_ROOT}/templates"

PROFILE_D_FILE="/etc/profile.d/hardn-paths.sh"

case "$1" in
    configure)
        # cross your fingers....
        ensure_hardn_account
        ensure_hardn_paths

        # Ensure system group 'hardn' exists, then create system user 'hardn' in that group
        if ! getent group hardn >/dev/null; then
            addgroup --system hardn || true
        fi
        if ! getent passwd hardn >/dev/null; then
            adduser --system --ingroup hardn --home /var/lib/hardn --shell /usr/sbin/nologin hardn || true
        fi

        # If package was installed via sudo by a logged-in user, add them to the hardn group
        if [ -n "${SUDO_USER:-}" ] && [ "${SUDO_USER}" != "root" ] && getent passwd "${SUDO_USER}" >/dev/null; then
            usermod -aG hardn "${SUDO_USER}" 2>/dev/null || true
            echo "Added ${SUDO_USER} to group 'hardn'"
        fi

        # Create symlink for backward compatibility
        if [ -f /usr/bin/hardn ] && [ ! -e /usr/bin/hardn ]; then
            ln -sf /usr/bin/hardn /usr/bin/hardn
        fi

        # Runtime directories: ensure group 'hardn' can read/write and new files inherit the group
        mkdir -p /var/log/hardn /var/lib/hardn /var/lib/hardn/backups 2>/dev/null || true
        chown root:hardn /var/log/hardn /var/lib/hardn /var/lib/hardn/backups 2>/dev/null || true
        chmod 2770 /var/log/hardn /var/lib/hardn /var/lib/hardn/backups 2>/dev/null || true

        # Initialize Legion baseline database directory
        mkdir -p /var/lib/hardn/baselines 2>/dev/null || true
        chown root:hardn /var/lib/hardn/baselines 2>/dev/null || true
        chmod 2770 /var/lib/hardn/baselines 2>/dev/null || true

        # Config directory: readable by group but writable only by root (maintain security)
        if [ -d /etc/hardn ]; then
            chown root:hardn /etc/hardn 2>/dev/null || true
            chmod 750 /etc/hardn 2>/dev/null || true
        fi

        # Ensure packaged data under /usr/share/hardn is owned by root:hardn and inherits group
        if [ -d /usr/share/hardn ]; then
            chown -R root:hardn /usr/share/hardn 2>/dev/null || true
            # setgid so new files/dirs inherit the hardn group; world read allowed for static assets
            find /usr/share/hardn -type d -exec chmod 2755 {} + 2>/dev/null || true
        fi

        # Deploy default config if missing
        if [ ! -f /etc/hardn/hardn.conf ] && [ -f "${TEMPLATE_DIR}/hardn.conf" ]; then
            cp "${TEMPLATE_DIR}/hardn.conf" /etc/hardn/
            chown root:hardn /etc/hardn/hardn.conf
            chmod 640 /etc/hardn/hardn.conf
        fi

        # Ensure runtime dirs
        install -d -m 755 -o root -g root "${MODULE_DIR}" "${TOOL_DIR}" "${TEMPLATE_DIR}"

        # Sync legacy modules if present
        if [ -d "${LEGACY_MODULE_DIR}" ]; then
            cp -a "${LEGACY_MODULE_DIR}/." "${MODULE_DIR}/" || true
        fi

        # Sync legacy tools if present
        if [ -d "${LEGACY_TOOL_DIR}" ]; then
            cp -a "${LEGACY_TOOL_DIR}/." "${TOOL_DIR}/" || true
        fi

        # Symlink fallback if empty
        if [ -d "${LEGACY_MODULE_DIR}" ] && [ -z "$(ls -A "${MODULE_DIR}" 2>/dev/null)" ]; then
            rm -rf "${MODULE_DIR}" && ln -s "${LEGACY_MODULE_DIR}" "${MODULE_DIR}"
        fi
        if [ -d "${LEGACY_TOOL_DIR}" ] && [ -z "$(ls -A "${TOOL_DIR}" 2>/dev/null)" ]; then
            rm -rf "${TOOL_DIR}" && ln -s "${LEGACY_TOOL_DIR}" "${TOOL_DIR}"
        fi

        # Ensure scripts are executable
        chmod +x "${MODULE_DIR}"/*.sh 2>/dev/null || true
        chmod +x "${TOOL_DIR}"/*.sh 2>/dev/null || true
        chmod +x /usr/bin/hardn-service-manager 2>/dev/null || true

        # Create global env vars for module/tool search paths
        echo "#!/bin/sh
# HARDN module and tool search paths
export HARDN_MODULE_PATH=\"${MODULE_DIR}:${LEGACY_MODULE_DIR}\"
export HARDN_TOOL_PATH=\"${TOOL_DIR}:${LEGACY_TOOL_DIR}\"
" > "${PROFILE_D_FILE}"
        chmod 644 "${PROFILE_D_FILE}"

        # Install systemd service files if they exist
        if [ -d "systemd" ]; then
            for service_file in systemd/*.service; do
                if [ -f "$service_file" ]; then
                    install -m 644 "$service_file" /lib/systemd/system/ || true
                fi
            done
        fi
        
        # Reload and enable systemd services updated 
        if command -v systemctl >/dev/null 2>&1 && systemctl is-system-running >/dev/null 2>&1 || [ "${container:-}" != "docker" ]; then
            systemctl daemon-reload || true
            # Enable services in dependency order
            systemctl enable hardn.service || true
            systemctl enable legion-daemon.service || true
            systemctl enable hardn-api.service || true
            systemctl enable hardn-monitor.service || true
        else
            echo "systemd not available or not running, skipping service enablement."
        fi

        ICON_SRC="/usr/share/pixmaps/hardn-gui.jpeg"
        if [ -f "${ICON_SRC}" ]; then
            # Determine a non-root invoking user safely
            TARGET_USER="${SUDO_USER:-}"
            if [ -z "${TARGET_USER}" ] || [ "${TARGET_USER}" = "root" ]; then
                # logname may fail in some contexts; ensure it's valid
                _LOGNAME=$(logname 2>/dev/null || true)
                if [ -n "${_LOGNAME}" ] && getent passwd "${_LOGNAME}" >/dev/null; then
                    TARGET_USER="${_LOGNAME}"
                fi
            fi

            # Only deploy per-user files when we have a resolvable, non-root user
            if [ -n "${TARGET_USER}" ] && [ "${TARGET_USER}" != "root" ] && getent passwd "${TARGET_USER}" >/dev/null; then
                TARGET_HOME=$(getent passwd "${TARGET_USER}" | cut -d: -f6)
                if [ -n "${TARGET_HOME}" ] && [ -d "${TARGET_HOME}" ]; then
                    DESKTOP_DIR="${TARGET_HOME}/.local/share/applications"
                    ICON_DIR="${TARGET_HOME}/.local/share/icons"
                    DESKTOP_PATH="${DESKTOP_DIR}/hardn-gui.desktop"
                    ICON_DEST="${ICON_DIR}/hardn-gui.jpeg"

                    # Create directories as root, then copy files then chown carefully
                    install -d -m 755 "${DESKTOP_DIR}" "${ICON_DIR}" || true
                    install -m 644 "${ICON_SRC}" "${ICON_DEST}" || cp -f "${ICON_SRC}" "${ICON_DEST}" 2>/dev/null || true
                    cat <<EOF > "${DESKTOP_PATH}"
[Desktop Entry]
Type=Application
Name=HARDN Control Center
Comment=Launch the HARDN monitoring console
Exec=/usr/bin/hardn-gui
Icon=${ICON_DEST}
Terminal=false
Categories=Security;System;
StartupNotify=true
EOF
                    # Change ownership only if user exists and home directory matches
                    chown "${TARGET_USER}:${TARGET_USER}" "${DESKTOP_PATH}" "${ICON_DEST}" 2>/dev/null || true
                    echo "HARDN desktop launcher deployed for ${TARGET_USER}."
                else
                    echo "HARDN desktop launcher skipped: unable to resolve home for ${TARGET_USER}."
                fi
            else
                echo "HARDN desktop launcher skipped: no non-root sudo user detected."
            fi
        fi

        echo "HARDN installed successfully."
        echo "Run 'sudo hardn' to execute modules, or 'sudo hardn run-module <name>' for one module."
        echo "Environment paths set in ${PROFILE_D_FILE}"
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *)
        echo "postinst called with unknown argument '$1'" >&2
        exit 1
        ;;
esac

exit 0
