# debian/postinst
#!/bin/bash
set -e

LEGACY_MODULE_DIR="/usr/lib/hardn-xdr/src/setup/modules"
LEGACY_TOOL_DIR="/usr/lib/hardn-xdr/src/setup/tools"

RUNTIME_ROOT="/usr/share/hardn"
MODULE_DIR="${RUNTIME_ROOT}/modules"
TOOL_DIR="${RUNTIME_ROOT}/tools"
TEMPLATE_DIR="${RUNTIME_ROOT}/templates"

PROFILE_D_FILE="/etc/profile.d/hardn-paths.sh"

case "$1" in
    configure)

        # Create system user and group
        if ! getent passwd hardn >/dev/null; then
            adduser --system --group --home /var/lib/hardn --shell /usr/sbin/nologin hardn
        fi

        # Create symlink for backward compatibility
        if [ -f /usr/bin/hardn-xdr ] && [ ! -e /usr/bin/hardn ]; then
            ln -sf /usr/bin/hardn-xdr /usr/bin/hardn
        fi

        # Log directory
        install -d -m 750 -o hardn -g hardn /var/log/hardn

        # Config directory
        if [ -d /etc/hardn ]; then
            chown root:hardn /etc/hardn
            chmod 750 /etc/hardn
        fi

        # Backup directory
        install -d -m 750 -o hardn -g hardn /var/lib/hardn/backups

        # Deploy default config if missing
        if [ ! -f /etc/hardn/hardn.conf ] && [ -f "${TEMPLATE_DIR}/hardn.conf" ]; then
            cp "${TEMPLATE_DIR}/hardn.conf" /etc/hardn/
            chown root:hardn /etc/hardn/hardn.conf
            chmod 640 /etc/hardn/hardn.conf
        fi

        # Ensure runtime dirs
        install -d -m 755 -o root -g root "${MODULE_DIR}" "${TOOL_DIR}" "${TEMPLATE_DIR}"

        # Sync legacy modules if present
        if [ -d "${LEGACY_MODULE_DIR}" ]; then
            cp -a "${LEGACY_MODULE_DIR}/." "${MODULE_DIR}/" || true
        fi

        # Sync legacy tools if present
        if [ -d "${LEGACY_TOOL_DIR}" ]; then
            cp -a "${LEGACY_TOOL_DIR}/." "${TOOL_DIR}/" || true
        fi

        # Symlink fallback if empty
        if [ -d "${LEGACY_MODULE_DIR}" ] && [ -z "$(ls -A "${MODULE_DIR}" 2>/dev/null)" ]; then
            rm -rf "${MODULE_DIR}" && ln -s "${LEGACY_MODULE_DIR}" "${MODULE_DIR}"
        fi
        if [ -d "${LEGACY_TOOL_DIR}" ] && [ -z "$(ls -A "${TOOL_DIR}" 2>/dev/null)" ]; then
            rm -rf "${TOOL_DIR}" && ln -s "${LEGACY_TOOL_DIR}" "${TOOL_DIR}"
        fi

        # Ensure scripts are executable
        chmod +x "${MODULE_DIR}"/*.sh 2>/dev/null || true
        chmod +x "${TOOL_DIR}"/*.sh 2>/dev/null || true

        # Create global env vars for module/tool search paths
        echo "#!/bin/sh
# HARDN-XDR module and tool search paths
export HARDN_MODULE_PATH=\"${MODULE_DIR}:${LEGACY_MODULE_DIR}\"
export HARDN_TOOL_PATH=\"${TOOL_DIR}:${LEGACY_TOOL_DIR}\"
" > "${PROFILE_D_FILE}"
        chmod 644 "${PROFILE_D_FILE}"

        # Install systemd service files if they exist
        if [ -d "systemd" ]; then
            for service_file in systemd/*.service; do
                if [ -f "$service_file" ]; then
                    install -m 644 "$service_file" /lib/systemd/system/ || true
                fi
            done
        fi
        
        # Reload and enable systemd services
        if command -v systemctl >/dev/null; then
            systemctl daemon-reload
            # Enable services in dependency order
            systemctl enable hardn.service || true
            systemctl enable hardn-xdr.service || true
            # hardn-monitor is deprecated/optional - don't enable by default
            # systemctl enable hardn-monitor.service || true
        else
            echo "systemctl not found, skipping service enablement."
        fi

        echo "HARDN installed successfully."
        echo "Run 'sudo hardn-xdr' to execute modules, or 'sudo hardn-xdr run-module <name>' for one module."
        echo "Environment paths set in ${PROFILE_D_FILE}"
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

exit 0